# environments/loc/common/config/supervisor/queue-phone-check.conf
[program:QueuePhoneCheckJob]
process_name=%(program_name)s_%(process_num)02d
command={{ COMMON_CONFIG_SUPERVISOR_QUEUE_PHP }} {{ APP_PATH }}/yii queue-phone-check/listen
stdout_logfile={{ APP_CONSOLE_LOGFILE_PATH }}/queue-phone-check.log
autostart=true
autorestart=true
user={{ COMMON_CONFIG_SUPERVISOR_QUEUE_USER }}
stopsignal=KILL
numprocs={{ COMMON_CONFIG_SUPERVISOR_QUEUEPHONECHECK_NUMPROCS }}