# environments/loc/common/config/supervisor/socket-server.conf
[program:SoketServer]
process_name=%(program_name)s_%(process_num)02d
command={{ COMMON_CONFIG_SUPERVISOR_QUEUE_PHP }} {{ APP_PATH }}/yii websocket-server-v3/start
stdout_logfile={{ APP_CONSOLE_LOGFILE_PATH }}/socket-server.log
autostart=true
autorestart=true
user={{ COMMON_CONFIG_SUPERVISOR_QUEUE_USER }}
stopsignal=KILL
numprocs={{ COMMON_CONFIG_SUPERVISOR_QUEUESOCKETSERVER_NUMPROCS }}

[inet_http_server]
port=127.0.0.1:9001
username={{ COMMON_CONFIG_SUPERVISOR_SOCKETSERVER_USERNAME }}
password={{ COMMON_CONFIG_SUPERVISOR_SOCKETSERVER_PASSWORD }}