# environments/stage/common/config/supervisor/queue-system-services.conf
[program:QueueSystemServices]
process_name=%(program_name)s_%(process_num)02d
command={{ COMMON_CONFIG_SUPERVISOR_QUEUE_PHP }} {{ APP_PATH }}/yii queue-system-services/listen
stdout_logfile={{ APP_CONSOLE_LOGFILE_PATH }}/queue-system-services.log
autostart=true
autorestart=true
user={{ COMMON_CONFIG_SUPERVISOR_QUEUE_USER }}
stopsignal=KILL
numprocs={{ COMMON_CONFIG_SUPERVISOR_QUEUESYSTEMSERVICES_NUMPROCS }}